# deploy/data/etc/xdg/autostart/wesop.desktop
[Desktop Entry]
Type=Application
Name=WeSOP
Exec=/opt/wesop/run.sh
Name[en_US]=wesop
